# systemd/cynara.target
[Unit]
Description=cynara sockets
DefaultDependencies=true

