# systemd/cynara-admin.socket
[Socket]
ListenStream=/run/cynara/cynara-admin.socket
SocketMode=0600
SmackLabelIPIn=@
SmackLabelIPOut=@

Service=cynara.service

[Unit]
Wants=cynara.target
Before=cynara.target

[Install]
WantedBy=sockets.target
